# desktop/motor-server.service
[Unit]
Description=Raspberry Pi Motor Control Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/desktop
ExecStart=/usr/bin/python3 /home/pi/desktop/motor_server.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
